# deploy/systemd/codex-cli-over-telegram.service
[Unit]
Description=Codex CLI over Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gnu
Group=gnu
WorkingDirectory=/home/gnu/codex-cli-over-telegram
Environment=HOME=/home/gnu
EnvironmentFile=/etc/codex-cli-over-telegram/env
EnvironmentFile=-/etc/codex-cli-over-telegram/deploy.env
ExecStart=/usr/bin/node /home/gnu/codex-cli-over-telegram/dist/index.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/gnu /srv/dev

[Install]
WantedBy=multi-user.target
